# container/build.target
[Unit]
Description=The target for build.service
Requires=systemd-logind.service systemd-user-sessions.service dbus.service
